# meta-enapter-core/recipes-enapter/enapter-cron-hook/files/enapter-cron-hook.timer
[Unit]
Description=Enapter Cron Hook Timer Service

[Timer]
Unit=enapter-cron-hook.service
OnCalendar=*-*-* *:00/01:00

[Install]
WantedBy=timers.target
