# meta-enapter-core/recipes-core/rauc/files/99-rauc.rules
SUBSYSTEM=="block", KERNEL!="*part*", ENV{DEVTYPE}=="disk", ACTION=="add", RUN+="/usr/bin/enapter-link-data-disk %k"
