# meta-enapter-core/recipes-enapter/enapter-podman-prepare/files/enapter-podman-prepare.service
[Unit]
Description=Enapter Podman Prepare Service

After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/bin/enapter-podman-prepare
RemainAfterExit=yes
Type=oneshot
